# debian/yandex-portoinit.postinst
#!/bin/sh
# postinst script

##DEBHELPER##

case "$1" in
    configure)
        chown root:porto /usr/sbin/portoinit || exit 2
	invoke-rc.d yandex-portoinit reload || :
    ;;
esac

exit 0
