# debian/yandex-porto.postinst
#!/bin/sh
# postinst script

##DEBHELPER##

case "$1" in
    configure)
        chown root:porto /usr/sbin/portoctl || exit 2
	invoke-rc.d yandex-porto reload || :
    ;;
esac

exit 0
